# data/lomiri-location-service-trust-stored.service.in
[Unit]
Description=Location Services Trust Store Daemon
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/trust-stored-skeleton --remote-agent DBusRemoteAgent --bus=system --local-agent MirAgent --trusted-mir-socket=${XDG_RUNTIME_DIR}/mir_socket_trusted --for-service LomiriLocationService --with-text-domain lomiri-location-service --store-bus session

[Install]
WantedBy=graphical-session.target
