# data/ubuntu-location-service-trust-stored.service.in
[Unit]
Description=Location Services Trust Store Daemon
Requires=ubuntu-location-service.service
After=ubuntu-location-service.service

[Service]
Type=simple
Exec=@CMAKE_INSTALL_FULL_BINDIR@/trust-stored-skeleton --remote-agent DBusRemoteAgent --bus=system --local-agent MirAgent --trusted-mir-socket=/var/run/user/$(id -u)/mir_socket_trusted --for-service UbuntuLocationService --with-text-domain ubuntu-location-service --store-bus session

[Install]
WantedBy=multi-user.target
